# bot.conf
[program:GoodOptbot]
directory=/root/GoodOptbot
command=python3 app.py

autostart=true
autorestart=true
stderr_logfile=/root/GoodOptbot/logfile_err.log
stdout_logfile=/root/GoodOptbot/logfile.log